# backup.crontab
*/15 * * * * /bin/sh /home/yzhdong/mobvoi/wooz/WOOZ/dumpDb.sh

